# scripts/supervisord_basic_prod.conf
[supervisord]
logfile =  %(here)s/../storage/supervisor_log/supervisord_ke3_basic_prod.log
logfile_maxbytes = 10MB
logfile_backups=3
loglevel = debug
pidfile =  %(here)s/../storage/supervisor_log/supervisord_ke3_basic_prod.pid
nodaemon = false
minfds = 2048
umask = 022
identifier = supervisor
nocleanup = true
childlogdir =  %(here)s/../storage/supervisor_log
strip_ansi = false
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 3
stderr_logfile_maxbytes = 10MB
stderr_logfile_backups = 3
directory = %(here)s

[program:kirim_worker_1]
command=php app_worker.php -t 1
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_2]
command=php app_worker.php -t 2
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_3]
command=php app_worker.php -t 3
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_4]
command=php app_worker.php -t 4
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_5]
command=php app_worker.php -t 5
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_6]
command=php app_worker.php -t 6
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_7]
command=php app_worker.php -t 7
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_8]
command=php app_worker.php -t 8
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_9]
command=php app_worker.php -t 9
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_10]
command=php app_worker.php -t 10
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_11]
command=php app_worker.php -t 11
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_12]
command=php app_worker.php -t 12
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_13]
command=php app_worker.php -t 13
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_14]
command=php app_worker.php -t 14
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_15]
command=php app_worker.php -t 15
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_16]
command=php app_worker.php -t 16
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_17]
command=php app_worker.php -t 17
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_18]
command=php app_worker.php -t 18
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_19]
command=php app_worker.php -t 19
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8

[program:kirim_worker_20]
command=php app_worker.php -t 0
autostart=true
autorestart=true
startsecs=10
startretries=99
exitcodes=0,2
process_name = %(program_name)s_%(process_num)02d
numprocs = 8